# EV3MenuEv3dev/debian/openrobertalab.postinst
#!/bin/sh
# postinst script for openrobertalab

set -e

case "$1" in
  configure)
    # Ask the bus to reload the config file
    if [ -x "/etc/init.d/dbus" ]; then
      invoke-rc.d dbus force-reload || true
    fi
  ;;
  abort-upgrade|abort-remove|abort-deconfigure)
  ;;
  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
  ;;
esac

#DEBHELPER#
